# bybitbot.service
[Unit]
Description=Bybit Trading Bot Web Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bybitbotgpt
ExecStart=/usr/bin/python3 /root/bybitbotgpt/cli.py web
Restart=always

[Install]
WantedBy=multi-user.target
